# containers/scheduler_ingest/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[program:router]
directory=/app
process_name=%(program_name)s_%(process_num)02d
command=python -m containers.scheduler_ingest.router.main --config /app/containers/scheduler_ingest/config/ingest.yaml --router-id %(process_num)d
numprocs=16
numprocs_start=0
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:ingestor]
directory=/app
process_name=%(program_name)s_%(process_num)02d
command=python -m containers.scheduler_ingest.ingestor.main --config /app/containers/scheduler_ingest/config/ingest.yaml --ingestor-id %(process_num)d
numprocs=16
numprocs_start=0
; When GOLDEN_DISCOVERY_RANKER_V1_INGEST_INLINE_ENABLED=true each of these
; 16 ingestor processes runs the same GoldenDiscoveryRuntimeScorer inference
; pipeline as the background scorer container. The scorer fans the 4 heads
; out across a thread pool; without these caps BLAS / OpenMP would
; re-parallelise *under* each of those threads (16 ingestors x 4
; head-threads x BLAS_NUM_CPUs = quadratic context-switch thrash that the
; previous inline-scoring deploy hit). Pinning BLAS to one thread keeps the
; head fan-out as the only parallelism layer. The env block is set
; unconditionally so an operator flipping the inline flag does not also
; have to remember to restart with new env.
environment=PYTHONPATH="/app",OMP_NUM_THREADS="1",OPENBLAS_NUM_THREADS="1",MKL_NUM_THREADS="1"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:feature_extractor]
directory=/app
process_name=%(program_name)s_%(process_num)02d
command=python -m containers.scheduler_ingest.feature_extractor.main --config /app/containers/scheduler_ingest/config/ingest.yaml --extractor-id %(process_num)d
numprocs=16
numprocs_start=0
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:stats_aggregator]
directory=/app
command=python -m containers.scheduler_ingest.stats.main --config /app/containers/scheduler_ingest/config/ingest.yaml
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=20
stopsignal=TERM
stopwaitsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
